# heartbeat_monitor.service
[Unit]
Description=Heartbeat Monitor — Windows online detection & audio source switching
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=audio
WorkingDirectory=/home/audio
ExecStart=/usr/bin/python3 /home/audio/heartbeat_monitor.py
Restart=always
RestartSec=5

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=heartbeat-monitor

[Install]
WantedBy=default.target