# heartbeat_monitor.service
[Unit]
Description=Heartbeat Monitor — Windows online detection & audio source switching
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=guido
WorkingDirectory=/home/guido/driving-pi-system-scripts
ExecStart=/home/guido/driving-pi-system-scripts/start_heartbeat.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=heartbeat-monitor

[Install]
WantedBy=default.target
